# streaming/conf/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:srs]
command=/usr/local/srs/objs/srs -c /usr/local/srs/conf/srs.conf
directory=/usr/local/srs
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/usr/local/srs/logs/srs.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=SRS_LOG_LEVEL=%(ENV_SRS_LOG_LEVEL)s

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
directory=/etc/nginx
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface